# src/wm/config/vaultwm.service
[Unit]
Description=VaultOS Window Manager
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/local/bin/vaultwm
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical.target

