# src/utils/services/vaultos-theme-daemon.service
[Unit]
Description=VaultOS Theme Management Daemon
After=graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/vaultos-theme-daemon
Restart=on-failure
User=%i

[Install]
WantedBy=graphical.target

